# ansible/roles/vpnserver/files/reboot.cron
@reboot 180 && /sbin/iptables-restore < /etc/iptables.rules


#