# roles/nginx/templates/letsencrypt.cron
# distributed from ansible
PATH=/usr/bin:/bin:/usr/sbin
23 16 * * * /usr/local/bin/certbot renew -n
