# ansible/roles/afraidorg/templates/afraidorg.cron.j2
37 1,13 * * * /usr/bin/curl https://freedns.afraid.org/dynamic/update.php?{{ afraidorg_keys[inventory_hostname] }} >> /tmp/afraidorg_curl.log

#