# ansible/roles/nginx/templates/letsencrypt.cron
# distributed from ansible
PATH=/usr/bin:/bin:/usr/sbin
*/10 * * * * /usr/local/bin/certbot renew -n
